# whatsapp-tracker/whatsapp-tracker.service
# Systemd service file for the WhatsApp Progress Tracker
# Copy to: /etc/systemd/system/whatsapp-tracker.service
# Then run:
#   sudo systemctl daemon-reload
#   sudo systemctl enable whatsapp-tracker
#   sudo systemctl start whatsapp-tracker

[Unit]
Description=WhatsApp Progress Tracker
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/side-projects/whatsapp-tracker
EnvironmentFile=/path/to/side-projects/whatsapp-tracker/.env
ExecStart=/usr/bin/python3 run.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
